# tg-to-max.service
[Unit]
Description=TG-to-Max file bridge
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tg-to-max
EnvironmentFile=-/root/tg-to-max/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/tg-to-max/.venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
